# gswitch.service
[Unit]
Description=Automatic GPU switching service

[Service]
ExecStart=gswitch boot

[Install]
WantedBy=multi-user.target
